# container-assets/systemd/containerd.service.d/override.conf
[Service]
# Override default root to use persistent storage volume (mirrors docker --data-root /.data/docker)
ExecStart=
ExecStart=/usr/bin/containerd --root /.data/containerd
