# container-assets/systemd/docker.socket.d/override.conf
[Socket]
# Make the Docker socket world-accessible so the discobot user can use Docker
SocketMode=0666
